# ai-devops-monitor.timer
[Unit]
Description=AI DevOps Monitor Timer
After=ai-devops-monitor.service

[Timer]
# Run every 10 minutes
OnBootSec=1min
OnUnitActiveSec=10min
Unit=ai-devops-monitor.service

[Install]
WantedBy=timers.target
